# extras/androidtv/rootfs/etc/s6-overlay/s6-rc.d/create-user/up
#!/command/execlineb -P
# shellcheck disable=all

with-contenv
multisubstitute
{
  importas PUID PUID
  importas PGID PGID
}
/etc/s6-overlay/scripts/create-user
